# backend/systemd/archive-sync.service
[Unit]
Description=Process search data for the Character Archive
After=network.target

[Service]
User=server
Group=server
SyslogIdentifier=archive-sync
ExecStart=/srv/server/venv/bin/python /srv/server/chub-archive-server/backend/elastic-sync.py --service
LimitNOFILE=30000

[Install]
WantedBy=multi-user.target
